# supervisord-import.conf
[program:import]
command=/import.sh
numprocs=1
autostart=true
